# systemd/fips-dash.service
[Unit]
Description=FIPS Dash API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=fips
Group=fips
WorkingDirectory=/opt/fips-dash
EnvironmentFile=/opt/fips-dash/.env
ExecStart=/usr/local/bin/bun run src/server/index.ts
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
